# assets/crontab
* * * * * /home/ec2-user/report-metrics.sh
